# systemd/cool_down.service
[Unit]
Description=Reduce temp by 1 degree 

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo "-1" > /run/thermostat/offset.txt'


