# systemd/offset_reset.service
[Unit]
Description=Reset the Timer

[Service]
Type=oneshot
ExecStart=/bin/rm /run/thermostat/offset.txt


