# systemd/offset_reset.timer
[Unit]
Description=Reset the offset at 7AM

[Timer]
OnCalendar=*-*-* 07:00:00
Unit=set-offset@0.service

[Install]
WantedBy=thermostat.timer

