# systemd/offset_reset.service
[Unit]
Description=Reset the Temp offset to 0

[Service]
Type=oneshot
ExecStart=/bin/rm /run/thermostat/offset.txt


