# systemd/set-offset@.service
[Unit]
Description=Set thermostat offset to %I

[Service]
Type=oneshot
ExecStart=/bin/sh -c "echo %i > /run/thermostat/offset.txt"

